# backend/prodsupervisord.conf
[supervisord]
nodaemon=true

[program:django]
environment = PYTHONUNBUFFERED=1
command=sh -c "python -u manage.py flush --noinput && python manage.py makemigrations && python manage.py migrate && gunicorn -b 0.0.0.0:8000 netdis.wsgi"
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:celery]
command=celery -A netdis worker -l info
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0